# roles/graphdb/templates/graphdb@.service
[Unit]
Description=GraphDB with database directory %I

[Service]
Type=simple
WorkingDirectory=
ExecStart=/opt/graphdb/graphdb-se-{{ graphdb_version }}/bin/graphdb -s -Dgraphdb.home=%i