# roles/gstore/templates/gstore@.service
[Unit]
Description=gStore with database %i

[Service]
Type=simple
WorkingDirectory={{ database_base_dir }}/gstore/%i
ExecStart=/opt/gstore/bin/ghttp %i