# roles/fuseki/templates/fuseki@.service
[Unit]
Description=Starts a Fuseki server with database %i

[Service]
Type=simple
ExecStart=/opt/fuseki/apache-jena-fuseki-{{ fuseki_version }}/fuseki-server --loc={{ database_base_dir }}/fuseki/%i