# roles/iguana/templates/iguana-corecontroller.service
[Unit]
Description=The Iguana Corecontroller
BindsTo=rabbitmq-server.service
After=rabbitmq-server.service

[Service]
Type=simple
ExecStart=/usr/bin/java -cp /opt/iguana-{{ iguana_version.version_number }}/corecontroller.jar org.aksw.iguana.cc.controller.MainController /opt/iguana-{{ iguana_version.version_number }}/core.properties
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
