# roles/blazegraph/templates/blazegraph@.service
[Unit]
Description=Starts a Fuseki server with database %i

[Service]
Type=simple
User=blazegraph
WorkingDirectory={{ database_base_dir }}/blazegraph/%i
ExecStart=/usr/bin/java -jar /opt/blazegraph/blazegraph.jar 9999 %i /opt/blazegraph/application.properties
SuccessExitStatus=143