# roles/virtuoso/templates/virtuoso@.service
[Unit]
Description=A Virtuoso Database server using config %i

[Service]
Type=simple
ExecStart=/opt/virtuoso/{{ virtuoso_version }}/virtuoso-opensource/bin/virtuoso-t -c /opt/virtuoso/virtuoso-%i.ini +foreground