# roles/tentris/templates/tentris@.service
[Unit]
Description=Tentris server serving dataset %I

[Service]
ExecStart=/opt/tentris/tentris_server --logstdout -f %I
